# lib/systemd/system/NetworkManager.service.d/10-configure-for-ptusb0.conf
[Service]
ExecCondition=/usr/bin/raspi-config nonint is_pifour
ExecStartPre=/usr/lib/pt-display-port/ptusb0/setup/ptusb0-setup
ExecStartPost=/usr/lib/pt-display-port/dhcp-server/setup/dhcp-server-setup
Restart=on-failure
